# packaging/nfc-manager.service
[Unit]
Description=NFC manager

[Service]
Type=simple
ExecStart=/usr/bin/nfc-manager-daemon
KillSignal=SIGKILL

[Install]
WantedBy=multi-user.target
